# roles/shrink_logs/templates/shrink-logs.service.j2
[Unit]
Description=Shrink system logs to 100M

[Service]
Type=oneshot
ExecStart=/usr/bin/journalctl --vacuum-size=100M